# backtesting/deploy/orb-backtester.service
[Unit]
Description=ORB Backtester API
After=network-online.target main-db.service
Wants=network-online.target main-db.service

[Service]
Type=simple
WorkingDirectory=/opt/orb-backtester
EnvironmentFile=/opt/orb-backtester/.env
ExecStart=/usr/local/bin/uv run python scripts/run_server.py --host 0.0.0.0 --port 8200
Restart=always
RestartSec=10

StandardOutput=journal
StandardError=journal
SyslogIdentifier=orb-backtester

TimeoutStopSec=30
KillMode=mixed
LimitNOFILE=65536

[Install]
WantedBy=multi-user.target
